# LaunchAgent/dev.pruiz.LoginShot.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.pruiz.LoginShot</string>

    <key>ProgramArguments</key>
    <array>
        <string>__APP_EXECUTABLE_PATH__</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>__STDOUT_PATH__</string>

    <key>StandardErrorPath</key>
    <string>__STDERR_PATH__</string>
</dict>
</plist>
